# systemd/zomboid-watchdog.timer
[Unit]
Description=Zomboid B42 watchdog, every 3 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=3min
AccuracySec=30s

[Install]
WantedBy=timers.target
