# systemd/zomboid-b42.service
[Unit]
Description=Project Zomboid B42 Dedicated Server (box64/ARM emulation)
After=network-online.target zomboid-ciopfs.service
Wants=network-online.target
Requires=zomboid-ciopfs.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/zomboid-server
# box64 tuning — MUST match the [ProjectZomboid64] block in /etc/box64.box64rc.
# STRONGMEM=3 is critical: lower values deadlock the JVM during boot.
Environment=BOX64_DYNAREC_STRONGMEM=3
Environment=BOX64_DYNAREC_SAFEFLAGS=2
Environment=BOX64_DYNAREC_ALIGNED_ATOMICS=1
Environment=BOX64_JVM=0
Environment=BOX64_SSE42=0
Environment=BOX64_MAXCPU=4
Environment=BOX64_LOG=0
# NOTE: replace CHANGE_ME_ADMIN_PASSWORD with your admin password.
ExecStart=/bin/bash /opt/zomboid-server/start-server.sh -servername servertest -adminpassword CHANGE_ME_ADMIN_PASSWORD
# start-server.sh always `exit 0`, which masks JVM crashes -> `on-failure` would never fire.
# Use `always` so a crashed server still restarts.
Restart=always
RestartSec=20
TimeoutStopSec=120
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
