# templates/zomboid-watchdog.service
[Unit]
Description=Zomboid B42 boot-hang watchdog (box64)

[Service]
Type=oneshot
Environment=PZ_CONSOLE=__HOME__/Zomboid/server-console.txt
ExecStart=/usr/local/sbin/zomboid-watchdog.sh
